# bash/geesome-ipfs-restart.service
[Unit]
Description=Restart IPFS container

[Service]
Type=oneshot
ExecStart=/usr/bin/docker restart go_ipfs

[Install]
WantedBy=multi-user.target


